# _unused/logstash/start.conf
[program:logstash]
command=[% container.config.CONTAINER.PATHS.CONFIG %]/logstash.sh
priority = 3
startsecs=5
stdout_logfile=syslog
stderr_logfile=syslog

